# bbswitch.service
[Unit]
Description=Aryan-BBSwitch

[Service]
ExecStart=/usr/sbin/insmod /opt/modules/bbswitch.ko

[Install]
WantedBy=multi-user.target